# kindle/upstart/kindle-desk-card.conf
# /etc/upstart/kindle-desk-card.conf
# Autostart the Kindle desk-card poll loop after the framework is ready.
# Install:  mntroot rw && cp this file to /etc/upstart/ && mntroot ro && reboot
# Verify:   status kindle-desk-card

start on started framework
stop on stopping framework

respawn
respawn limit 10 60

script
  exec sh /mnt/us/extensions/desk-card/poll.sh >> /tmp/poll.log 2>&1
end script
